# src/airflow-2.0/supervisor/rsync.conf
[program:rsync]
command                 = bash "/appz/scripts/rsync.sh"
autostart               = false
autorestart             = false
stdout_logfile  = /appz/log/%(program_name)s.log
stderr_logfile  = /appz/log/%(program_name)s.log
