# src/mariadb-10.11/supervisor/fluentd.conf
[program:fluentd]
command                 = bash /appz/scripts/fluentd_run.sh
autostart               = true
autorestart             = true
stdout_logfile  = /appz/log/%(program_name)s.log
stderr_logfile  = /appz/log/%(program_name)s.log
startretries	= 10
redirect_stderr = true
