# src/airflowsch-2.0/supervisor/airflowsch.conf
[program:airflow_scheduler]
command=bash -c "/appz/scripts/airflow.sh scheduler"
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/appz/log/%(program_name)s.log
stderr_logfile=/appz/log/%(program_name)s.error.log
autorestart=true